# lbdb-fetch-maildir.service
[Unit]
Description=Update lbdb from Maildir

[Service]
Type=simple
ExecStart=%h/.neomutt/bin/lbdb-fetch-maildir.sh "%h/.mail"